# autostart/Entry.desktop
[Desktop Entry]
Exec=sh /usr/bin/MyApp.sh
Terminal=false

